# wc.service
[Unit]
Description = WC

[Service]
ExecStart = /home/pi/wc/wc.py
Restart = always
Type = simple

[Install]
WantedBy = multi-user.target
